# system_files/bluefin/etc/profile.d/open.sh
# shellcheck shell=bash
alias open="xdg-open &>/dev/null"
