# system_files/shared/usr/share/ublue-os/dx/quadlets/cockpit-dx.container
# DX-Next: Cockpit web console in Podman (optional menu component).
[Unit]
Description=Cockpit Web Service (DX-Next)
After=network-online.target
Wants=network-online.target sshd.service

[Container]
Image=quay.io/cockpit/ws:latest
ContainerName=cockpit-ws
Volume=/:/host:rslave
PodmanArgs=--privileged --pid=host

Network=host

[Install]
WantedBy=multi-user.target
