# system_files/amd/usr/share/ublue-os/amd-stacks/pytorch-lab/pytorch-lab.container
[Unit]
Description=PyTorch Lab — ROCm ${ROCM_VERSION} (JupyterLab workspace :8888)
After=network-online.target

[Container]
ContainerName=pytorch-lab
Image=docker.io/rocm/pytorch:rocm${ROCM_VERSION}_ubuntu24.04_py3.12_pytorch_release_2.10.0
AddDevice=/dev/kfd
AddDevice=/dev/dri
PodmanArgs=--security-opt=label=disable
PodmanArgs=--group-add=video
PodmanArgs=--ipc=host
Ulimit=memlock=-1:-1
Ulimit=stack=67108864:67108864
Network=pytorch-lab-network.network
PublishPort=8888:8888
Volume=%h/ai-workspaces/pytorch-lab:/workspace:z
Exec=jupyter lab --no-browser --ip=0.0.0.0 --port=8888 --allow-root \
     --ServerApp.token='' --ServerApp.password='' --notebook-dir=/workspace

[Service]
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
