# system_files/amd/usr/share/ublue-os/amd-stacks/ollama/ollama.container
[Unit]
Description=Ollama — ROCm inference server (Ollama + OpenAI-compatible API :11434)
After=network-online.target

[Container]
ContainerName=ollama
Image=docker.io/ollama/ollama:rocm
AddDevice=/dev/kfd
AddDevice=/dev/dri
PodmanArgs=--security-opt=label=disable
PodmanArgs=--group-add=video
Network=ollama-network.network
PublishPort=11434:11434
Volume=%h/ai-workspaces/ollama:/root/.ollama:z
Environment=OLLAMA_HOST=0.0.0.0

[Service]
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
