# system_files/amd/usr/share/ublue-os/amd-stacks/lemonade/lemonade.container
[Unit]
Description=Lemonade — AMD-native LLM server (ROCm · Vulkan · NPU · port 13305)
After=network-online.target

[Container]
ContainerName=lemonade
Image=ghcr.io/lemonade-sdk/lemonade-server:v10.7.0
AddDevice=/dev/kfd
AddDevice=/dev/dri
PodmanArgs=--security-opt=label=disable
PodmanArgs=--group-add=video
Network=lemonade-network.network
PublishPort=13305:13305
Volume=%h/ai-workspaces/lemonade/cache:/root/.cache/huggingface:z
Volume=%h/ai-workspaces/lemonade/llama:/opt/lemonade/llama:z
Volume=%h/ai-workspaces/lemonade/config:/root/.cache/lemonade:z

[Service]
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
