# system_files/nvidia/usr/share/ublue-os/nvidia-stacks/nim-llama3/nim-llama3.container
[Unit]
Description=NIM — Llama 3.1 8B Instruct (OpenAI-compatible API)
After=network-online.target

[Container]
ContainerName=nim-llama3
Image=nvcr.io/nim/meta/llama-3.1-8b-instruct:latest
AddDevice=nvidia.com/gpu=all
PodmanArgs=--security-opt=label=disable
PodmanArgs=--ipc=host
Ulimit=memlock=-1:-1
Ulimit=stack=67108864:67108864
Secret=ngc-api-key,type=env,target=NGC_API_KEY
Network=nim-llama3-network.network
PublishPort=8000:8000
Volume=%h/.cache/nim:/root/.cache/nim:z

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=600

[Install]
WantedBy=default.target
