# system_files/amd/usr/share/ublue-os/amd-stacks/llama-strix/llama-strix.container
[Unit]
Description=Llama Strix — Strix Halo-optimized llama.cpp server (gfx1151 · ROCm 7.2.4 · :8080)
After=network-online.target

[Container]
ContainerName=llama-strix
Image=docker.io/kyuz0/amd-strix-halo-toolboxes:rocm-7.2.4
AddDevice=/dev/kfd
AddDevice=/dev/dri
PodmanArgs=--security-opt=label=disable
PodmanArgs=--security-opt=seccomp=unconfined
PodmanArgs=--group-add=video
PodmanArgs=--group-add=render
Network=llama-strix-network.network
PublishPort=8080:8080
Volume=%h/ai-workspaces/llama-strix/models:/workspace/models:z
# -fa 1 and --no-mmap are REQUIRED on Strix Halo — omitting either causes crashes/slowdowns
# -ngl 999 offloads all layers to the iGPU (unified memory = no overflow penalty)
Exec=llama-server -m ${LLAMA_MODEL} -c 8192 -ngl 999 -fa 1 --no-mmap --host 0.0.0.0 --port 8080

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=default.target
